# crontab
30,45 3-17 * * * uv run main.py
0,15 4-18 * * * uv run main.py
